# deploy/systemd/mpquic@mp1.service.d/bd1-route.conf
[Service]
ExecStartPost=/bin/sh -c "sleep 1 && ip route replace default dev mp1 table bd1 && ip route replace 172.16.17.0/30 dev enp6s20.17 table bd1"
