# deploy/systemd/mpquic-mgmt.service
[Unit]
Description=MPQUIC Management API
Documentation=https://github.com/peppechiapparo/mpquic
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mpquic-mgmt \
    --listen 10.10.11.100:8080 \
    --instance-dir /etc/mpquic/instances
EnvironmentFile=-/etc/mpquic/mgmt.env
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/mpquic/instances
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mpquic-mgmt

[Install]
WantedBy=multi-user.target
