# deploy/systemd/mpquic-watchdog.timer
[Unit]
Description=Run MPQUIC tunnel watchdog every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=mpquic-watchdog.service

[Install]
WantedBy=timers.target
