# deploy/systemd/mpquic-watchdog.service
[Unit]
Description=MPQUIC tunnel watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/lib/mpquic/mpquic-tunnel-watchdog.sh
